# minimization/system/sysinit.target
[Unit]
Description=Empty sysinit target
