# dependencies/five-services-with-wants-requires-and-conflicts/system/a.service
[Unit]
Requires=b.service
After=b.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=sleep 1
ExecStop=sleep 1
