# dependencies/five-services-with-wants-requires-and-conflicts/system/b.service
[Unit]
Wants=c.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=sleep 1
ExecStart=false
ExecStop=sleep 1
