# dependencies/five-services-with-wants-requires-and-conflicts/system/halt.target
[Unit]
Description=A minimal halt target
Requires=halt.service
