# dependencies/six-services-with-ordering-and-requirement-dependencies/system/c.service
[Unit]
Before=e.service
Wants=a.service b.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=sleep 1
ExecStop=sleep 1
