# dependencies/six-services-with-ordering-and-requirement-dependencies/system/d.service
[Unit]
Before=f.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=sleep 5
ExecStop=sleep 1
