# dependencies/six-services-with-ordering-and-requirement-dependencies/system/default.target
[Unit]
Description=A minimal default target
Requires=systemd-journald.service sysinit.target
