# dependencies/six-services-with-ordering-and-requirement-dependencies/system/halt.service
[Unit]
Description=Halt systemd
DefaultDependencies=no

[Service]
ExecStart=systemctl --force halt
