# dependencies/six-services-with-ordering-dependencies/system/a.service
[Unit]
Before=c.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=sleep 3
ExecStop=sleep 1
