# dependencies/six-services-with-ordering-dependencies/system/c.service
[Unit]
Before=e.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=sleep 1
ExecStop=sleep 1
