# dependencies/six-services-with-ordering-dependencies/system/e.service
[Unit]
Before=f.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=sleep 1
ExecStop=sleep 1
