# dependencies/six-services-with-ordering-dependencies/system/systemd-journald.service
[Unit]
Description=Journal Service
DefaultDependencies=no

[Service]
ExecStart=/lib/systemd/systemd-journald
Sockets=systemd-journald.socket
